# ops/host-deploy.timer
[Unit]
Description=jekko host auto-deploy timer (poll jeryu every 2 minutes)
Requires=host-deploy.service

[Timer]
# Fire 30s after activation, then every 2 minutes.
OnActiveSec=30s
OnUnitActiveSec=2min
Unit=host-deploy.service
Persistent=true

[Install]
WantedBy=timers.target
